# user-cli/home/profile.d/rust.sh

# rust
if [ -d /opt/rustc ]; then
    export RUST_SRC_PATH="/opt/rustc/src"
fi
