# user-cli/home/profile.d/rust.sh
# RUST
if [ -d /opt/rustc ]; then
    #export PATH="$HOME/.cargo/RustDT/bin:$PATH"
    export RUST_SRC_PATH="/opt/rustc/src"
fi
